# hue-panel.service
[Unit]
Description=Hue Vacation Panel
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/hue-panel
ExecStart=/usr/bin/node /opt/hue-panel/.output/server/index.mjs
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=HUE_STORAGE_PATH=/opt/hue-panel/data
Environment=PORT=3000
Environment=HOST=0.0.0.0

[Install]
WantedBy=multi-user.target
